# home/.config/systemd/user/check-kernel-version.service
[Unit]
Description=Check kernel version every hour

[Service]
Type=oneshot
ExecStart=check-kernel-version notify
