# home/.config/systemd/user/check-kernel-version.timer
[Unit]
Description=Check kernel version every hour

[Timer]
OnCalendar=hourly
Unit=%i.service

[Install]
WantedBy=timers.target
